# distribution/systemd/jarvis-driver-manager.service
[Unit]
Description=JARVIS OS Driver Manager
Documentation=https://github.com/jarvis-os
After=jarvis-hal.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/lib/jarvis/jarvis-driver-manager
Restart=on-failure
RestartSec=10

[Install]
WantedBy=jarvis.target
