# distribution/systemd/jarvis-ui.service
[Unit]
Description=JARVIS OS Web UI
Documentation=https://github.com/jarvis-os
After=network.target

[Service]
Type=simple
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvis/frontend/web
ExecStart=/usr/bin/python -m http.server 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=jarvis.target
